# owasp-zap-weekly.service
[Unit]
Description=OWASP ZAP weekly build service

[Service]
Type=oneshot
ExecStart=/bin/bash "$HOME/.local/bin/owasp-zap-weekly.sh"
