# zap_weekly.timer
[Unit]
Description=ZAP weekly-build download timer

[Timer]
OnCalendar=Mon 21:00
Persistent=true

[Install]
WantedBy=timers.target
